# saltstack/inex/states/cumulus/etc/systemd/system/ntp.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/lib/ntp/ntp-systemd-wrapper
